# profile.d/dvm-env.sh
#dvm for everybody who uses it
#2018 tiagobrait

DVM_DIR="$HOME/.dvm/scripts"

[[ -s "$DVM_DIR/dvm" ]] && source "$DVM_DIR/dvm"
