# profile.d/flutter-env.sh
#!/bin/bash
#2018 tiaogbrait
#flutter env stuff

FLUT_DIR="$HOME/flutter"

[[ -d "${FLUT_DIR}/bin" ]] && export PATH="${FLUT_DIR}/bin":${PATH}
