# profile.d/jdk.sh
#! /bin/bash
#export JAVA_HOME to the path setup by archlinux-java

export JAVA_HOME="/usr/lib/jvm/default"
